# etc/systemd/system/chrony-wait-stop.path
[Unit]
Description=Monitor for timesync synchronization file

[Path]
PathExists=/run/systemd/timesync/synchronized

[Install]
WantedBy=multi-user.target
