# etc/systemd/system/copy-authorized-keys.path
[Path]
PathChanged=/boot/firmware/authorized_keys

[Install]
RequiredBy=network.target
