# etc/systemd/system/huaweicheck.service
[Unit]
Description=Check Huawei cellular connection
After=network.target

[Service]
Nice=10
ExecStart=/usr/local/sbin/huaweicheck
Type=oneshot
