# etc/systemd/system/tsconfig.service
[Unit]
Description=Web-based configuration and system management interface for trackIT Systems sensor stations
After=network.target

[Service]
User=pi
Group=pi
Environment="TERM=xterm-256color"
WorkingDirectory=/home/pi/tsconfig
ExecStart=python3 -m uvicorn app.main:app

Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
