# etc/systemd/system/wittypid-watcher.path
[Path]
PathChanged=/boot/firmware/schedule.yml

[Install]
WantedBy=multi-user.target
