# etc/systemd/system/wittypi-watcher.path
[Path]
PathChanged=/boot/firmware/schedule.wpi

[Install]
WantedBy=multi-user.target
