# etc/systemd/system/wittypid.service
[Unit]
Description=Power Management

[Service]
ExecStart=/usr/local/bin/wittypid --schedule /boot/firmware/schedule.yml -v

[Install]
WantedBy=multi-user.target
